# sys_config/lib/systemd/system/radio.service
# Unit file for kradio
[Unit]
Description=KRadio
After=network.target

[Service]
Type=simple
User=kradio
SupplementaryGroups=i2c
WorkingDirectory=/opt/kradio
ExecStart=/opt/kradio/.venv/bin/kradio

[Install]
WantedBy=multi-user.target
